# data/gnulib/profile.d/Makefile.am
# Makefile.am

include $(top_srcdir)/rules.mk

# data/script variables
gnulib = gnulib.sh

# data/scripts generated by configure (source file automatically distributed)
conf_data = $(gnulib)

profile_DATA = $(conf_data)

DISTCLEANFILES = $(conf_data)
