# data/bash/common/profile.d/android.sh
#! /bin/bash

pathmung $HOME/usr/share/android/tools after
