# data/bash/common/profile.d/arch-bin.sh
#! /bin/bash

# Set PATH so it includes user's private architecture-specific bin if it exists
pathmung $HOME/usr/$(/bin/arch)/bin
