# data/bash/common/profile.d/bash.sh
#! /bin/bash

# History Options
# ###############

# Don't put duplicate lines in the history.
HISTCONTROL="ignoredups"

# Ignore some controlling instructions
HISTIGNORE="[   ]*:&:bg:fg:exit"
