# data/bash/common/profile.d/conan.sh.in
# gnulib settings

export CONAN_CLANG_VERSIONS=@CONAN_CLANG_VERSIONS@
export CONAN_GCC_VERSIONS=@CONAN_GCC_VERSIONS@
export CONAN_PIP_COMMAND='@PYTHON@ -m pip'
export CONAN_USERNAME=@CONAN_USERNAME@
