# data/bash/common/profile.d/conan.sh.in
# gnulib settings

export CONAN_USERNAME=@CONAN_USERNAME@
export CONAN_PIP_COMMAND=@CONAN_PIP_COMMAND@
