# data/bash/common/profile.d/editor.sh.in
# editor settings

export EDITOR=@EDITOR@

export VISUAL=$EDITOR
export GIT_EDITOR=$EDITOR
export SVN_EDITOR=$EDITOR
