# data/bash/common/profile.d/editor.sh.in
#! /bin/bash

EDITOR=@EDITOR@

GIT_EDITOR=${EDITOR}

SVN_EDITOR=${EDITOR}
