# data/bash/common/profile.d/eldk.sh
#! /bin/bash

pathmung /opt/eldk/bin after
pathmung /opt/eldk/usr/bin after

CROSS_COMPILE="ppc_4xx-"
