# data/bash/common/profile.d/git.sh
export GIT_PS1_SHOWDIRTYSTATE=true
export GIT_PS1_SHOWSTASHSTATE=true
export GIT_PS1_SHOWUNTRACKEDFILES=true
export GIT_PS1_SHOWUPSTREAM='auto verbose name'
export GIT_PS1_DESCRIBE_STYLE='describe'
