# data/bash/common/profile.d/gnulib.sh.in
#! /bin/bash

GNULIB_SRCDIR=@GNULIB_SRCDIR@
