# data/bash/common/profile.d/gnulib.sh.in
# gnulib settings

export GNULIB_REFDIR=@GNULIB_REFDIR@
