# data/bash/common/profile.d/gnulib.sh.in
# gnulib settings

export GNULIB_SRCDIR=@GNULIB_SRCDIR@
