# data/bash/common/profile.d/go.sh.in
# go settings

export GOPATH=@GOPATH@
export GOBIN=@bindir@
