# data/bash/common/profile.d/gvim.sh.in
export GVIMINIT='let $MYGVIMRC="@vimdir@/gvimrc" | source $MYGVIMRC'
