# data/bash/common/profile.d/hadoop.sh.in
# hadoop settings

export HADOOP_HOME=@HADOOP_HOME@

pathmung @HADOOP_HOME@/bin after
