# data/bash/common/profile.d/hadoop.sh.in
HADOOP_HOME=@HADOOP_HOME@

pathmung @HADOOP_HOME@ after
