# data/bash/common/profile.d/hbase.sh.in
HBASE_HOME=@HBASE_HOME@

pathmung @HBASE_HOME@ after
