# data/bash/common/profile.d/idea.sh.in
IDEA_HOME=@IDEA_HOME@

pathmung @IDEA_HOME@/bin after
