# data/bash/common/profile.d/macosx.sh
# Mac OS X environment settings

# prevent tar from adding mac's extended attribute file to archives
export COPYFILE_DISABLE=true
