# data/bash/common/profile.d/scl.sh.in
# login to scl-enabled bash

# only do this in interactive shells
if echo $- | grep i &> /dev/null; then
    for scl in @SCL_LIST@; do
        if ! @SCL_ENABLED@ $scl; then
            source @SCL_SOURCE@ enable $scl
        fi
    done
fi
