# data/bash/common/profile.d/zz-path.sh.in
# Set PATH so it includes user's private bin dirs if they exist
pathmung @bindir@
pathmung @sbindir@
