# data/ruby/profile.d/bundle.sh.in
# bundle settings

export BUNDLE_USER_CONFIG=@BUNDLE_USER_CONFIG@
export BUNDLE_USER_CACHE=@BUNDLE_USER_CACHE@
export BUNDLE_USER_PLUGIN=@BUNDLE_USER_PLUGIN@
