# data/bash/profile.d/bash-completion.sh
# Completion options
# ##################

# These completion tuning parameters change the default behavior of 
# bash_completion:

# Define to access remotely checked-out files over passwordless ssh for CVS
export COMP_CVS_REMOTE=1

# Define to avoid flattening internal contents of tar files
export COMP_TAR_INTERNAL_PATHS=1
