# data/bash/profile.d/conan.sh.in
# gnulib settings

export CONAN_CLANG_VERSIONS=@CONAN_CLANG_VERSIONS@
export CONAN_GCC_VERSIONS=@CONAN_GCC_VERSIONS@
export CONAN_PIP_COMMAND=@CONAN_PIP_COMMAND@
export CONAN_USER_HOME=@CONAN_USER_HOME@
export CONAN_USERNAME=@CONAN_USERNAME@
