# data/bash/profile.d/vim.sh.in
export VIMINIT='let $MYVIMRC="@vimdir@/vimrc" | source $MYVIMRC'
