# data/bash/umbc/profile.d/Makefile.am
# Makefile.am

dist_profile_DATA = umbc.sh

if UMBC_ORACLE
dist_profile_DATA += umbc-oracle.sh
endif

if UMBC_CADENCE
dist_profile_DATA += umbc-cadence.sh
endif
