# data/bash/umbc/profile.d/umbc.sh
#! /bin/bash

pathmung /afs/umbc.edu/users/r/h/rheingan/pub/341 after
