# data/bash/umbc/profile.d/umbc-cadence.sh
#! /bin/bash

pathmung \
  /afs/umbc.edu/software/cadence/software_2008/IUS61/tools.lnx86/inca/bin after

CDS_LIC_FILE="5280@license1.umbc.edu"
LM_LICENSE_FILE="/afs/umbc.edu/software/cadence/etc/license.dat"
