# data/bash/umbc/profile.d/umbc-oracle.sh
#! /bin/bash

source /usr/local/oracle/oracle.profile
