# data/readline/profile.d/readline.sh.in
# readline settings
# vim: ft=sh :

export INPUTRC=@readlineconfdir@/inputrc
