# data/git/profile.d/git.sh.in
export GIT_AUTHOR_EMAIL='@GIT_AUTHOR_EMAIL@'
export GIT_AUTHOR_NAME='@GIT_AUTHOR_NAME@'
export GIT_PS1_SHOWDIRTYSTATE=yes
export GIT_PS1_SHOWSTASHSTATE=yes
export GIT_PS1_SHOWUNTRACKEDFILES=yes
export GIT_PS1_SHOWUPSTREAM='auto verbose name'
export GIT_PS1_SHOWCONFLICTSTATE=yes
export GIT_PS1_SHOWCOLORHINTS=yes
export GIT_PS1_DESCRIBE_STYLE='describe'
