# data/vim/profile.d/editor.sh.in
# editor settings

export EDITOR=@VIM@
export VISUAL=$EDITOR
export GIT_EDITOR=$EDITOR
