# data/vim/profile.d/vim.sh.in
export VIMINIT='let $MYVIMRC="@vimconfdir@/vimrc" | source $MYVIMRC'
