# data/gnupg/profile.d/Makefile.am
# Makefile.am

include $(top_srcdir)/rules.mk

# data/script variables
gpg = gpg.sh

# static data/scripts (not generated; require explicit distribution)
static_data = $(gpg)

dist_profile_DATA = $(static_data)
