# data/info/profile.d/info.sh.in
# man settings

export INFOPATH='@infodir@:'
