# data/go/profile.d/go.sh.in
# go settings

export GOPATH=@GOPATH@
export GOMODCACHE=@GOMODCACHE@
export GOBIN=@bindir@
